# tests/Command/cron.d/cronjobs-1
*/2 * * * * www-data /usr/bin/test -d
* * * * * root /usr/bin/mail -s "spam"
* 4 * FEB-DEC SUN-SAT root /usr/bin/mail -s "spam"
